# lp_hedge_backtest/viznago_api.service
[Unit]
Description=VIZNAGO FURY FastAPI Service
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/dev/trading/lp_hedge_backtest
EnvironmentFile=/var/www/dev/trading/lp_hedge_backtest/api/.env
ExecStart=/var/www/dev/trading/lp_hedge_backtest/venv/bin/uvicorn \
    api.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 1 \
    --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=viznago-api

[Install]
WantedBy=multi-user.target
